# startup_autorun.service
[Unit]
Description=startup script

[Service]
ExecStart=/home/debian/rover/startup.sh

[Install]
WantedBy=multi-user.target 
